# bridge/sbc/openautolink-eth-ssh.service
[Unit]
Description=OpenAutoLink Ethernet SSH Access (DHCP server on eth0)
After=network-pre.target
Before=network.target
Wants=network-pre.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/openautolink/setup-eth-ssh.sh start
ExecStop=/opt/openautolink/setup-eth-ssh.sh stop

[Install]
WantedBy=multi-user.target
