# bridge/sbc/openautolink-bt.service
[Unit]
Description=OpenAutoLink Bluetooth (BLE + HSP + AA RFCOMM)
After=bluetooth.target
Wants=bluetooth.service

[Service]
Type=simple
# Wait for BT hardware to appear (any SBC — vendor init services create hci0 at
# varying times). bluetooth.service has ConditionPathIsDirectory=/sys/class/bluetooth
# so it gets skipped on first systemd pass; udev triggers it once hci0 appears.
# This loop waits up to 30s for that to happen before starting aa_bt_all.py.
ExecStartPre=/bin/sh -c 'for i in $(seq 1 60); do [ -d /sys/class/bluetooth/hci0 ] && exit 0; sleep 0.5; done; echo "[bt] WARNING: hci0 not found after 30s, starting anyway"; exit 0'
EnvironmentFile=-/etc/openautolink.env
EnvironmentFile=-/boot/firmware/openautolink.env
ExecStart=/usr/bin/python3 -OO /opt/openautolink/scripts/aa_bt_all.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
