# bridge/sbc/openautolink-wireless.service
[Unit]
Description=OpenAutoLink Wireless (WiFi AP)
After=openautolink-network.service network.target bluetooth.target
Before=openautolink.service
Wants=bluetooth.service
# WiFi AP MUST start after bluetooth.target. If the AP comes up first, a
# non-default phone with cached WiFi credentials can TCP-connect before the
# BT RFCOMM gate and startup grace are armed, bypassing multi-phone safety.

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/openautolink.env
EnvironmentFile=-/boot/firmware/openautolink.env
# Wait for a wireless interface to appear (any SBC — vendor WiFi services load
# the driver module at varying times). Checks /sys/class/net/*/wireless which
# exists for any WiFi NIC regardless of driver or naming convention.
ExecStartPre=/bin/sh -c 'for i in $(seq 1 60); do for p in /sys/class/net/*/wireless; do [ -d "$p" ] && exit 0; done; sleep 0.5; done; echo "[wifi] WARNING: no wireless interface after 30s"; exit 0'
ExecStart=/opt/openautolink/start-wireless.sh
ExecStop=/opt/openautolink/stop-wireless.sh

[Install]
WantedBy=multi-user.target
