# bridge/sbc/openautolink-network.service
[Unit]
Description=OpenAutoLink Network Setup (car + SSH)
DefaultDependencies=no
Before=openautolink.service network.target
After=local-fs.target sysinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/openautolink.env
EnvironmentFile=-/boot/firmware/openautolink.env
ExecStart=/opt/openautolink/setup-network.sh setup
ExecStop=/opt/openautolink/setup-network.sh teardown

[Install]
WantedBy=multi-user.target
