# bridge/sbc/openautolink-wireless.service
[Unit]
Description=OpenAutoLink Wireless (WiFi AP)
After=openautolink-network.service network.target bluetooth.target
Before=openautolink.service
Wants=bluetooth.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/openautolink.env
EnvironmentFile=-/boot/firmware/openautolink.env
ExecStart=/opt/openautolink/start-wireless.sh
ExecStop=/opt/openautolink/stop-wireless.sh

[Install]
WantedBy=multi-user.target
